# niftywall.service
[Unit]
Description=NiftyWall Firewall Dashboard
After=network.target nftables.service

[Service]
User=root
Group=root
WorkingDirectory=/opt/niftywall
Environment="PATH=/opt/niftywall/venv/bin"
EnvironmentFile=/opt/niftywall/.env
ExecStart=/opt/niftywall/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8080 --timeout-keep-alive 120
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target